# natlas-server/deployment/natlas-server.service
[Unit]
Description=Natlas Server
After=network.target
 
[Service]
Type=simple
User=natlas
WorkingDirectory=/opt/natlas/natlas-server
ExecStart=/bin/bash /opt/natlas/natlas-server/run-server.sh
Restart=always

[Install]
WantedBy=multi-user.target